# exam/roles/haproxy/templates/haproxy_exporter.service
[Unit]
Description=HAProxy Prometheus exporter

[Service]
ExecStart=/usr/local/bin/haproxy_exporter

[Install]
WantedBy=multi-user.target
